# deploy/launchd/com.zora.ainews.fetch.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.zora.ainews.fetch</string>

  <!-- 执行本地抓取脚本 -->
  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>/Users/zora/CodeBuddy/ai-news-aggregator/scripts/local-fetch.sh</string>
  </array>

  <!-- 白天每 2 小时跑一次（10/12/14/16/18/20 点的第 5 分钟，错开云端整点调度）。
       电脑这些时刻处于关机/睡眠则该次跳过；下次到点正常执行。 -->
  <key>StartCalendarInterval</key>
  <array>
    <dict><key>Hour</key><integer>10</integer><key>Minute</key><integer>5</integer></dict>
    <dict><key>Hour</key><integer>12</integer><key>Minute</key><integer>5</integer></dict>
    <dict><key>Hour</key><integer>14</integer><key>Minute</key><integer>5</integer></dict>
    <dict><key>Hour</key><integer>16</integer><key>Minute</key><integer>5</integer></dict>
    <dict><key>Hour</key><integer>18</integer><key>Minute</key><integer>5</integer></dict>
    <dict><key>Hour</key><integer>20</integer><key>Minute</key><integer>5</integer></dict>
  </array>

  <key>StandardOutPath</key>
  <string>/Users/zora/CodeBuddy/ai-news-aggregator/logs/launchd.out.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/zora/CodeBuddy/ai-news-aggregator/logs/launchd.err.log</string>

  <!-- 加载时不立即执行，避免安装即触发；按计划时间运行 -->
  <key>RunAtLoad</key>
  <false/>
</dict>
</plist>
